# bsp-target-files/etc/cron.d/dynainstruments
# MIN HOUR DAY MONTH DAYOFWEEK  COMMAND
@hourly ID=save_swclock /etc/init.d/swclock --quiet save_now
